# devopin-agent.service
[Unit]
Description=Devopin System Monitoring Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=devopin
Group=devopin
WorkingDirectory=/opt/devopin-agent
ExecStart=/opt/devopin-agent/devopin-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devopin-agent

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/devopin-agent /var/log/devopin /etc/devopin /run/devopin-agent /var/lib/devopin-agent

[Install]
WantedBy=multi-user.target